# supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/local/bin/entrypoint.sh
